# PillowModule/ATD_LP_DSE_v2.2.1-main/ATD_LP_DSE_v2.2.1-main/services/atDSE.service
\[Unit]
Description= Atlab Service
After=network.target
StartLimitIntervalSec=0


[Service]
User=root
WorkingDirectory=/atd/src
ExecStart=/usr/bin/python3 main.py

[Install]
WantedBy=multi-user.target