# PillowModule/ATD_LP_DSE_v2.2.1-main/ATD_LP_DSE_v2.2.1-main/services/atOpenVPN.service
\[Unit]
Description= Atlab Service
After=network.target
StartLimitIntervalSec=0


[Service]
User=root
WorkingDirectory=/atd/src/func/opnvpn
ExecStart=/usr/bin/python3 opnvpn.py

[Install]
WantedBy=multi-user.target